# crontab
0 0 * * * node /root/nodemail/main.js >> /root/nodemail/crontab.log
